# supervisord.conf
[supervisord]

[program:websockets-HiveMind]
command = python server.py > /srv/http/HiveMindServer/logs/server.log
process_name = %(program_name)s_%(process_num)02d
numprocs = 1
autorestart = false
